# infra/systemd/trademetrix-web.service
[Unit]
Description=TradeMetrix Terminal — Next.js Frontend
After=network.target trademetrix-api.service
Wants=trademetrix-api.service

[Service]
Type=simple
User=nodejs
Group=nodejs
WorkingDirectory=/opt/trademetrix/apps/web
EnvironmentFile=/opt/trademetrix/apps/web/.env
ExecStart=/usr/bin/node server.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trademetrix-web
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
MemoryLimit=512M
CPUQuota=100%

[Install]
WantedBy=multi-user.target
