# infra/systemd/trademetrix-cron.timer
[Unit]
Description=TradeMetrix Terminal — Daily Bhavcopy Sync Trigger

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
